# supervisor.conf
[supervisord]
nodaemon=true

[program:cosmicpi-daq]
command=cosmicpi start
